# dot-config_systemd_user/nightynight@.service
[Unit]
Description=Set night mode %i

[Service]
Type=oneshot
ExecStart=/home/furios/.local/bin/nightynight %i
