# systemd_user/nightynight@.service
[Unit]
Description=Set night mode %i

[Service]
Type=oneshot
ExecStart=/usr/bin/nightynight %i
